# wol-nrw-fetch-levels.service
[Service]
Type=oneshot
ExecStart=%h/.local/bin/wol-nrw-fetch-levels.sh
