# wol-nrw-fetch-levels.timer
[Timer]
OnCalendar=*-*-* *:10:00

[Install]
WantedBy=default.target
